# local-vendors/aws-code-deploy/assets/aws/default/supervisor/confd.conf
[program:confd]
command=confd -interval 600 -backend dynamodb -table configuration-$SERVER_TYPE
numprocs=1
autostart=true
autorestart=true
priority=2
stderr_logfile=/var/log/confd.err.log
stdout_logfile=/var/log/confd.out.log
environment=AWS_REGION=eu-west-1
